# launchagents/com.broodlink.hugo-status.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.broodlink.hugo-status</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOMEBREW_PREFIX__/bin/hugo</string>
        <string>server</string>
        <string>--source</string>
        <string>__BROOD_DIR__/status-site</string>
        <string>--bind</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>1313</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>BROODLINK_CONFIG</key>
        <string>__BROOD_DIR__/config.toml</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>__BROOD_DIR__</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/hugo-status.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/hugo-status.stderr.log</string>
</dict>
</plist>
