# launchagents/com.broodlink.mcp-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.broodlink.mcp-server</string>
    <key>ProgramArguments</key>
    <array>
        <string>__BROOD_DIR__/target/release/mcp-server</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>BROODLINK_CONFIG</key>
        <string>__BROOD_DIR__/config.toml</string>
        <key>BROODLINK_STATUS_API_KEY</key>
        <string>__STATUS_API_KEY__</string>
        <key>PATH</key>
        <string>__HOMEBREW_PREFIX__/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>__BROOD_DIR__</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/mcp-server.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/mcp-server.stderr.log</string>
</dict>
</plist>
